# snort/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:snort]
command=snort -D
stdout_logfile=/var/log/supervisor.log
stderr_logfile=/var/log/supervisor.log
autorestart=true